# roles/homeassistant/files/systemd/homeassistant-restart.service
[Unit]
Description=Restart homeassistant

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart homeassistant.service
